# tools/systemd/manpkid.service
[Unit]
Description=ManPKI Main Daemon

[Service]
Type=simple
ExecStart=/usr/bin/manpkid

[Install]
WantedBy=multi-user.target
