# tools/upstart/manpkid.conf
description "ManPKI Main Daemon"
author "manpki@ferez.fr"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

exec /usr/sbin/manpkid
